# systemd/db.socket
[Unit]
Description=gunicorn socket

[Socket]
ListenStream=/srv/aclarknet-database/run/db/socket

[Install]
WantedBy=sockets.target
