# systemd/db.socket
[Unit]
Description=gunicorn socket

[Socket]
ListenStream=/run/db/socket

[Install]
WantedBy=sockets.target
